# ansible/roles/metrics/files/podman-influxdb.service
# container-621abb3eb01401fb30e9d0e795af9ad7ab3047a726509fc3e351cae44b4a3849.service
# autogenerated by Podman 3.1.2
# Fri May 14 20:47:21 UTC 2021

[Unit]
Description=Podman container-621abb3eb01401fb30e9d0e795af9ad7ab3047a726509fc3e351cae44b4a3849.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target
RequiresMountsFor=/var/lib/containers/storage /run/containers/storage

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start 621abb3eb01401fb30e9d0e795af9ad7ab3047a726509fc3e351cae44b4a3849
ExecStop=/usr/bin/podman stop -t 10 621abb3eb01401fb30e9d0e795af9ad7ab3047a726509fc3e351cae44b4a3849
ExecStopPost=/usr/bin/podman stop -t 10 621abb3eb01401fb30e9d0e795af9ad7ab3047a726509fc3e351cae44b4a3849
PIDFile=/run/containers/storage/overlay-containers/621abb3eb01401fb30e9d0e795af9ad7ab3047a726509fc3e351cae44b4a3849/userdata/conmon.pid
Type=forking

[Install]
WantedBy=multi-user.target default.target
