# ansible/roles/metrics/files/telegraf.service
# container-6e39b27d20c760c127abdb169bec445ca22b9ba1db7352050d01f89c52556e80.service
# autogenerated by Podman 3.1.2
# Fri May 14 21:28:25 UTC 2021

[Unit]
Description=Podman container-6e39b27d20c760c127abdb169bec445ca22b9ba1db7352050d01f89c52556e80.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target
RequiresMountsFor=/var/lib/containers/storage /run/containers/storage

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start 6e39b27d20c760c127abdb169bec445ca22b9ba1db7352050d01f89c52556e80
ExecStop=/usr/bin/podman stop -t 10 6e39b27d20c760c127abdb169bec445ca22b9ba1db7352050d01f89c52556e80
ExecStopPost=/usr/bin/podman stop -t 10 6e39b27d20c760c127abdb169bec445ca22b9ba1db7352050d01f89c52556e80
PIDFile=/run/containers/storage/overlay-containers/6e39b27d20c760c127abdb169bec445ca22b9ba1db7352050d01f89c52556e80/userdata/conmon.pid
Type=forking

[Install]
WantedBy=multi-user.target default.target
