# ansible/podman-grafana.service
# container-5599d1ba608cd996780450806817f5f45c6753fff2ffe6576cf640c9eef70f80.service
# autogenerated by Podman 3.4.2
# Sat Jul 16 21:31:33 UTC 2022

[Unit]
Description=Podman container-5599d1ba608cd996780450806817f5f45c6753fff2ffe6576cf640c9eef70f80.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/run/containers/storage

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start 5599d1ba608cd996780450806817f5f45c6753fff2ffe6576cf640c9eef70f80
ExecStop=/usr/bin/podman stop -t 10 5599d1ba608cd996780450806817f5f45c6753fff2ffe6576cf640c9eef70f80
ExecStopPost=/usr/bin/podman stop -t 10 5599d1ba608cd996780450806817f5f45c6753fff2ffe6576cf640c9eef70f80
PIDFile=/run/containers/storage/overlay-containers/5599d1ba608cd996780450806817f5f45c6753fff2ffe6576cf640c9eef70f80/userdata/conmon.pid
Type=forking

[Install]
WantedBy=multi-user.target default.target
