# ansible/podman-influxdb.service
# container-5f59c82cd92108c2eca01e4dfd19db9e522f7a397a6313783f438dff212ab56c.service
# autogenerated by Podman 3.4.2
# Sun Nov 28 11:35:01 UTC 2021

[Unit]
Description=Podman container-5f59c82cd92108c2eca01e4dfd19db9e522f7a397a6313783f438dff212ab56c.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/run/containers/storage

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start 5f59c82cd92108c2eca01e4dfd19db9e522f7a397a6313783f438dff212ab56c
ExecStop=/usr/bin/podman stop -t 10 5f59c82cd92108c2eca01e4dfd19db9e522f7a397a6313783f438dff212ab56c
ExecStopPost=/usr/bin/podman stop -t 10 5f59c82cd92108c2eca01e4dfd19db9e522f7a397a6313783f438dff212ab56c
PIDFile=/run/containers/storage/overlay-containers/5f59c82cd92108c2eca01e4dfd19db9e522f7a397a6313783f438dff212ab56c/userdata/conmon.pid
Type=forking

[Install]
WantedBy=multi-user.target default.target
